# deploy/systemd/greenblock-serial-bridge.service
[Unit]
Description=GreenBlock Arduino Serial Bridge
After=network-online.target greenblock-backend.service
Wants=network-online.target
Requires=greenblock-backend.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Greenblock/greenblock-backend
EnvironmentFile=/home/pi/Greenblock/greenblock-backend/.env
ExecStart=/home/pi/Greenblock/greenblock-backend/.venv/bin/python serial_bridge.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
